# deploy/templates/cf-realip-sync.service.j2
[Unit]
Description=Update Cloudflare real_ip ranges and reload nginx
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cf-realip-sync
