# deploy/templates/cf-realip-sync.timer.j2
[Unit]
Description=Run cf-realip-sync weekly and at boot catch-up

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h
Unit=cf-realip-sync.service

[Install]
WantedBy=timers.target
