# deploy/templates/cf-ufw-sync.timer.j2
[Unit]
Description=Run cf-ufw-sync daily

[Timer]
OnCalendar=02:15
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
